# files/etc/profile.d/rbenv.sh
export RBENV_ROOT="/usr/lib/rbenv"
[ -d "$RBENV_ROOT" ] && eval "$(rbenv init -)"
